# scripts/systemd_risd.service
[Unit]
Description=Remote Input Service Daemon

[Service]
ExecStart=/usr/local/bin/risd

[Install]
WantedBy=multi-user.target
